# src/http-server.service
Unit]
Description=HTTP Server for Kerberos operations via Web form

[Service]
ExecStart=/usr/local/bin/http-server
