# src/init-django.service
[Unit]
Description=Populate Django database and admin user
Before=httpd.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/init-django
